# caddy/rootfs/etc/s6/php/finish
#!/bin/bash
/bin/s6-svscanctl -t /etc/s6
