# caddy/rootfs/etc/s6/php/run
#!/bin/bash

if [ -f ./setup ]
then
  source ./setup
fi

if [ -f ./custom ]
then
  source ./custom
fi

echo "Starting php-fpm!"
/usr/bin/php-fpm
