# caddy/rootfs/etc/s6/caddy/run
#!/bin/bash

if [ -f ./setup ]
then
  source ./setup
fi

if [ -f ./custom ]
then
  source ./custom
fi

/usr/bin/gosu \
  caddy \
  /usr/bin/caddy \
  -conf=/etc/caddy/caddyfile
