# nginx/rootfs/etc/s6/nginx/run
#!/bin/bash

if [ -f ./setup ]
then
  source ./setup
fi

/usr/sbin/nginx
