# nginx/rootfs/etc/s6/nginx/run
#!/bin/sh

if [ -f ./setup ]
then
  ./setup
fi

/usr/sbin/nginx
