# apache/rootfs/etc/s6/apache/run
#!/bin/sh

if [ -f ./setup ]
then
  ./setup
fi

if [ -f ./custom ]
then
  ./custom
fi

/usr/sbin/httpd -DFOREGROUND
