# nginx/rootfs/etc/s6/php/run
#!/bin/sh

if [ -f ./setup ]
then
  source ./setup
fi

if [ -f ./custom ]
then
  source ./custom
fi

/usr/bin/php-fpm
