# nginx/rootfs/etc/s6/php/run
#!/bin/bash

if [ -f ./setup ]
then
  source ./setup
fi

if [ -f ./custom ]
then
  source ./custom
fi

for FILE in custom.d/*.sh
do
  source ${FILE}
done

exec /usr/sbin/php-fpm7
